# config/wsgi_connect_to_db.conf
description	"WSGI connect to db"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/bin/python /srv/load_testing/test_cases/wsgi/db/connect/wsgi_connect_to_db.py
