# config/wsgi_read_from_disk.conf
description	"Simple WSGI"

start on filesystem or runlevel [2345]
stop on runlevel [!2345]

respawn

exec /usr/bin/python /srv/load_testing/test_cases/wsgi/disk_read/wsgi_read_from_disk.py
